# apache2-smf/apache2.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    Apache2 SMF Manifest
    Chris Miles 20050624
    http://chrismiles.info/systemsadmin/solaris/articles/setting-up-smf-for-apache-and-zope-on-solaris-10
-->

<service_bundle type='manifest' name='site:apache2'>

<service
        name='site/apache2'
        type='service'
        version='1'>

        <create_default_instance enabled='false' />

        <single_instance />

        <dependency name='paths'
            grouping='require_all'
            restart_on='error'
            type='path'>
                <service_fmri value='file://localhost/usr/apache2/2.2/bin/apachectl' />
                <service_fmri value='file://localhost/usr/apache2/2.2/bin/httpd' />
                <service_fmri value='file://localhost/etc/apache2/2.2/httpd.conf' />
        </dependency>

        <dependency name='network'
            grouping='require_any'
            restart_on='error'
            type='service'>
                <service_fmri value='svc:/network/service' />
        </dependency>

        <dependent
                name='apache2_multi-user'
                grouping='optional_all'
                restart_on='none'>
                <service_fmri value='svc:/milestone/multi-user' />
        </dependent>

        <!--
                The timeout needs to be large enough to wait for startup.
        -->
        <exec_method
            type='method'
            name='start'
            exec='/lib/svc/method/apache2 start'
            timeout_seconds='60' />

        <exec_method
            type='method'
            name='stop'
            exec='/lib/svc/method/apache2 stop'
            timeout_seconds='60' />

        <stability value='Unstable' />

        <template>
                <common_name>
                        <loctext xml:lang='C'>
                        Apache 2 (httpd)
                        </loctext>
                </common_name>
                <documentation>
                        <manpage title='httpd' section='8' />
                </documentation>
        </template>
</service>

</service_bundle>
